# installer/systemd/mirror-frontend.service
[Unit]
Description=Smart Mirror Frontend (SvelteKit, Node)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/mirror/config.env
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=3000
Environment=NODE_EXTRA_CA_CERTS=/etc/mirror/ha.crt
# WorkingDirectory rewritten by installer to match actual repo path
WorkingDirectory=/opt/mirror/frontend
ExecStart=/usr/bin/node build/index.js
Restart=always
RestartSec=5
User=mirror
Group=mirror

# systemd-managed log dir — created automatically under /var/log/mirror
LogsDirectory=mirror
LogsDirectoryMode=0755

# Hardening (kept lenient so it actually starts on a fresh box).
# NoNewPrivileges must be OFF: the admin endpoint shells out to
# `sudo -n /usr/bin/systemctl reboot` for the HA "Reboot PC" control,
# and sudo refuses to elevate when this flag is on.
NoNewPrivileges=false
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true

[Install]
WantedBy=multi-user.target
